# configs/apps/satellite/root_overlay/etc/profile.d/02-rhasspy-satellite.sh
export PATH=/mnt/persist/rhasspy/bin:$PATH
